# deployments/cloud/educates/lab-modern-spring-data/workshop/profile.d/01-rmq-operator-post-setup.sh
export RABBITMQ_SERVER_NAME="rabbitmq"
export RMQ_KUB_NAMESPACE="rabbitmq-system"